# packaging/debian_amd64/lib/systemd/system/fio-bco-import@.path
[Unit]
Description=openbank FIO Bank Connection tenant %i

[Path]
PathChanged=/etc/init/fio-bco.conf

[Install]
WantedBy=fio-bco-import@.service
