# packaging/debian/fio-bco-watcher.path
[Unit]
Description=openbank FIO Bank connection configuration watcher
After=fio-bco.service

[Path]
PathChanged=/etc/fio-bco/conf.d
Unit=fio-bco-watcher.service

[Install]
WantedBy=paths.target
