# packaging/debian/fio-bco.path
[Unit]
Description=openbank FIO Bank Connection services namespace
PartOf=fio-bco.service
ConditionPathExists=/etc/fio-bco/conf.d/init.conf

[Path]
PathChanged=/etc/fio-bco/conf.d/init.conf
PathModified=/etc/fio-bco/conf.d/init.conf
Unit=fio-bco.service

[Install]
WantedBy=fio-bco.service
