# packaging/debian_armhf/lib/systemd/system/fio-bco-rest.path
[Unit]
Description=openbank FIO Bank Connection REST service

[Path]
PathChanged=/etc/init/fio-bco.conf

[Install]
WantedBy=fio-bco-rest.service
