# packaging/debian_arm64/lib/systemd/system/fio-bco.path
[Unit]
Description=openbank FIO Bank Connection services namespace
PartOf=fio-bco.service
ConditionPathExists=/etc/init/fio-bco.conf

[Path]
PathChanged=/etc/init/fio-bco.conf
PathModified=/etc/init/vault.conf
Unit=fio-bco.service

[Install]
WantedBy=fio-bco.service
